# new_game/supervisord.conf

[program:router]
directory=/Users/john/dev/go/src/github.com/jnovikov/hackforces/new_game
command=router/router
autorestart=true
stdout_logfile=./router.log
stderr_logfile=./router.log


[inet_http_server]
port = 127.0.0.1:9999
username = admin
password = admin123


[supervisord]
logfile=./supervis/log/supervisord.log ; supervisord log file
logfile_maxbytes=50MB       ; maximum size of logfile before rotation
logfile_backups=10          ; number of backed up logfiles
loglevel=info               ; info, debug, warn, trace
pidfile=./supervis/run/supervisord.pid ; pidfile location
nodaemon=false              ; run supervisord as a daemon
minfds=1024                 ; number of startup file descriptors
minprocs=200                ; number of process descriptors


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use unix:// schem for a unix sockets.