# scripts/templates/supervisord.conf
{% for router in routers %}
[program:{{ router }}]
directory=%(here)s
command=%(here)s/{{ router }}
stdout_logfile=%(here)s/%(program_name)s.log
autorestart=true
redirect_stderr=true
{% endfor %}

[unix_http_server]
file=supervisor_server.sock  ; path to your socket file
username = admin
password = adminkeklol

[inet_http_server]
port = 127.0.0.1:8999
username = admin
password = adminkeklol


[supervisord]
logfile=%(here)s/supervisord.log ; supervisord log file
logfile_maxbytes=50MB       ; maximum size of logfile before rotation
logfile_backups=10          ; number of backed up logfiles
loglevel=info               ; info, debug, warn, trace
pidfile=%(here)s/supervisord.pid ; pidfile location
nodaemon=false              ; run supervisord as a daemon
minfds=1024                 ; number of startup file descriptors
minprocs=200                ; number of process descriptors


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://supervisor_server.sock ; use unix:// schem for a unix sockets.

[include]
files = service_controller*/*.conf
