# new_game/service_2/test.conf
[program:flag_handler-2]
directory=.
command=flag_handler/flag_handler
stdout_logfile=supervis/log/fh.log
stderr_logfile=supervis/log/fh.log
autorestart=true

[program:flag_adder-2]
directory=.
command=flag_adder/flag_adder
autorestart=true
stdout_logfile=supervis/log/fa.log
stderr_logfile=supervis/log/fa.log

[program:init_handler-2]
directory=.
command=init_handler/init_handler
autorestart=true
stdout_logfile=supervis/log/ih.log
stderr_logfile=supervis/log/ih.log